# app/crontab
0 * * * * /usr/bin/freshclam --quiet
# <- This is a newline character